# automation/com.wmprojekt.refresh-mvp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.wmprojekt.refresh-mvp</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/steffengorsdorf/Projects/wm-projekt/scripts/refresh_mvp.command</string>
		<string>--auto-publish</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/steffengorsdorf/Projects/wm-projekt</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Hour</key>
			<integer>5</integer>
			<key>Minute</key>
			<integer>5</integer>
			<key>Month</key>
			<integer>7</integer>
			<key>Day</key>
			<integer>19</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>3</integer>
			<key>Minute</key>
			<integer>10</integer>
			<key>Month</key>
			<integer>7</integer>
			<key>Day</key>
			<integer>20</integer>
		</dict>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardOutPath</key>
	<string>/Users/steffengorsdorf/Library/Logs/wm-projekt/refresh_mvp.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/steffengorsdorf/Library/Logs/wm-projekt/refresh_mvp.err.log</string>
</dict>
</plist>
